# hostex-chat-backend.service
[Unit]
Description=Hostex Chat Backend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/hostex-chat
EnvironmentFile=/path/to/hostex-chat/.env
ExecStart=/usr/bin/node backend/index.js
Restart=always
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
